# configs/controller/setup-00-init/systemd/var-mnt-data.automount
[Unit]
Description=Automount External drive

[Automount]
Where=/var/mnt/data

[Install]
WantedBy=multi-user.target
